# build-system/config/includes.chroot/etc/xdg/autostart/ridos-firstboot.desktop
[Desktop Entry]
# RIDOS OS v1.1 — First Boot Setup Wizard Autostart
# Copyright (C) 2026 RIDOS OS Project — GPL v3
# Runs setup_wizard.py on first login if not yet configured

Type=Application
Name=RIDOS First Boot Setup
Comment=Run RIDOS OS setup wizard on first boot
Exec=bash -c "[ ! -f /etc/ridos/.setup_complete ] && xfce4-terminal --title='RIDOS OS Setup' --geometry=90x40 --command='sudo python3 /opt/ridos/bin/setup_wizard.py' || true"
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=3
